# 参考文件/GitHub精品参考/miclaw_api_bridge/openwrt/files/etc/init.d/miclaw_api_bridge
#!/bin/sh /etc/rc.common
# procd service for miclaw_api_bridge headless server on OpenWrt.

START=95
STOP=10
USE_PROCD=1

PROG=/usr/bin/miclaw_api_bridge
CONF=miclaw_api_bridge

# Persist session/settings here (survives reboots; NOT cleared on sysupgrade
# only if added to /etc/sysupgrade.conf — see README).
DATA_DIR=/etc/miclaw_api_bridge

start_service() {
	# Defaults; overridable via /etc/config/miclaw_api_bridge
	local enabled host port verbose tls tls_cert tls_key

	config_load "$CONF"
	config_get_bool enabled  main enabled 1
	config_get      host     main host    "0.0.0.0"
	config_get      port     main port    "8765"
	config_get_bool verbose  main verbose 0
	config_get_bool tls      main tls     0
	config_get      tls_cert main tls_cert ""
	config_get      tls_key  main tls_key  ""

	[ "$enabled" = "1" ] || {
		echo "miclaw_api_bridge disabled in /etc/config/$CONF" >&2
		return 1
	}

	mkdir -p "$DATA_DIR/config" "$DATA_DIR/data"

	procd_open_instance
	procd_set_param command "$PROG" server --host "$host" --port "$port"

	# HTTPS: --tls auto-generates a self-signed cert under $DATA_DIR/config;
	# optionally point to your own PEM cert/key.
	[ "$tls" = "1" ] && procd_append_param command --tls
	[ -n "$tls_cert" ] && procd_append_param command --tls-cert "$tls_cert"
	[ -n "$tls_key" ] && procd_append_param command --tls-key "$tls_key"

	# No OS keyring on OpenWrt → store session on disk.
	procd_set_param env \
		MICLAW_API_BRIDGE_DISABLE_KEYRING=1 \
		XDG_CONFIG_HOME="$DATA_DIR/config" \
		XDG_DATA_HOME="$DATA_DIR/data" \
		RUST_LOG="miclaw_api_bridge_lib=$( [ "$verbose" = 1 ] && echo debug || echo info )"

	procd_set_param respawn 3600 5 0
	procd_set_param stdout 1
	procd_set_param stderr 1
	procd_close_instance
}

service_triggers() {
	procd_add_reload_trigger "$CONF"
}

reload_service() {
	stop
	start
}
